# s6/rootfs/etc/s6-overlay/s6-rc.d/ecowitt2mqtt/run
#!/command/with-contenv sh
ecowitt2mqtt
